# ansible/host-config/roles/custom_configs/templates/bmc-fan-metrics.timer.j2
[Unit]
Description=Poll BMC fan RPMs every 10 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
