# ansible/host-config/roles/custom_configs/templates/bmc-fan-metrics.service.j2
[Unit]
Description=BMC fan RPM metrics poller for prox0
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bmc-fan-metrics.py
EnvironmentFile=/etc/hdd-fan-curve/credentials
StandardOutput=journal
StandardError=journal
