# ansible/host-config/roles/custom_configs/templates/disk-temp-exporter.service.j2
[Unit]
Description=Export USB-SSD SMART temperatures to node-exporter textfile collector
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disk-temp-exporter.py
StandardOutput=journal
StandardError=journal
